# datasets/dataset16/etc/cron.d/docsearch-cleanup
0 4 * * 0 root /opt/docsearch/scripts/cleanup.sh
