# datasets/dataset16/etc/cron.d/docsearch-backup
0 2 * * * root /opt/docsearch/scripts/backup.sh
